# scripts/ingest_options.timer
[Unit]
Description=Timer for option contracts ingestion
Requires=ingest_options.service

[Timer]
# Start the first run at 21:00 America/Los_Angeles today (or next boot if missed),
# then repeat every 3 days after each run.
OnCalendar=21:00 America/Los_Angeles
Persistent=true
AccuracySec=1s
OnUnitActiveSec=3d

[Install]
WantedBy=timers.target
