# scripts/ingest_snapshots.timer
[Unit]
Description=Timer for option snapshots ingestion
Requires=ingest_snapshots.service

[Timer]
# Run daily at 23:30 PST (11:30 PM Pacific Time)
# Note: systemd uses UTC, PST is UTC-8, so 23:30 PST = 07:30 UTC next day
OnCalendar=daily
OnCalendar=*-*-* 07:30:00
Persistent=true

# Also allow manual triggering
AccuracySec=1s

[Install]
WantedBy=timers.target
