# scripts/ingest_options.service
[Unit]
Description=Ingest Option Contracts from Polygon API
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/jingyi/PycharmProjects/strategy-tester

Environment="PATH=/home/jingyi/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Load service-specific .env with ingest-options OTEL attributes
EnvironmentFile=/home/jingyi/PycharmProjects/strategy-tester/.env.options

# Increase file descriptor limit for Prisma connection pool
LimitNOFILE=65536

# Run the ingest_options zsh script
ExecStart=/home/jingyi/PycharmProjects/strategy-tester/scripts/ingest_options.zsh

# Consider timeout (exit code 124) as success
SuccessExitStatus=124

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ingest_options

[Install]
WantedBy=default.target
